# ansible/backup/templates/cleanup_backups_daily.timer
[Unit]
Description=Run backup cleanup daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
