# ansible/backup/templates/pg-backup-my-wiki.service
[Unit]
Description=Logical PostgreSQL Backup of my_wiki Database

[Service]
Type=oneshot
User=postgres
ExecStart=/var/lib/postgresql/backup.sh
